# 48-block-mouse.rules
SUBSYSTEM=="usb", ATTR{bDeviceClass}=="00", \
    ATTRS{bInterfaceClass}=="03", \
    RUN+="/bin/sh -c 'echo 0 > /sys$env{DEVPATH%/*}/authorized'"
