# 48-block-mouse.rules
SUBSYSTEM=="usb", ENV{DEVTYPE}=="usb_interface", \
    ATTR{bInterfaceClass}=="03", \
    RUN+="/usr/bin/block-hid.sh $env{DEVPATH}"
